# webservice/deployment/crontab
*/5 * * * * php /var/www/html/assets/cronjob_mail.php